# nixos/.config/autostart/jottad.desktop
[Desktop Entry]
Name=Jottad
Exec=systemd-run --user --unit=jottad jottad
